# WebGoat/webgoat.conf
[supervisord]
nodaemon=true

[program:tomcat8]
command=/opt/apache-tomcat-8.5.9/bin/catalina.sh start
