# bWAPP/bwapp.conf
[supervisord]
nodaemon=true

[program:mysqld]
command=/etc/init.d/mysql restart

[program:apache2]
command=/etc/init.d/apache2 restart
