# DVWA-WooYun-edition/dvwa_wooyun.conf
[supervisord]
nodaemon=true

[program:mysql]
command=service mysql restart

[program:apache2]
command=service apache2 restart
