# WAVSEP/wavsep.conf
[supervisord]
nodaemon=true

[program:mysqld]
command=/etc/init.d/mysql restart

[program:tomcat8]
command=/opt/apache-tomcat-8.5.11/bin/catalina.sh start
