# OWASP Security Shepherd/security-shepherd.conf
[supervisord]
nodaemon=true

[program:mongod]
command=/usr/bin/mongod

[program:mysqld]
command=/usr/bin/mysqld_safe

[program:tomcat7]
command=service tomcat7 start
